# roles/os_settings/templates/rdp-start.timer.j2
[Unit]
Description=Start {{ rdp_service }} at {{ rdp_timer_start }}

[Timer]
OnCalendar={{ rdp_timer_start }}
Persistent=true
RandomizedDelaySec=0

[Install]
WantedBy=timers.target
